# skills/superlooper/skill/templates/launchd.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- superlooper keep-alive template (plan Task 10). The DEFAULT way to run the loop is a
     visible `superlooper run` in a cmux tab; this launchd job is the keep-alive option for
     unattended stretches. Placeholders ({label}, {superlooper_bin}, {repo_path}, {state_home})
     are substituted literally. Label + log paths use the same slash-safe <owner>__<repo>
     normalization as the state home. The external-watchdog contract is the state home's
     state/runner.heartbeat (epoch, stamped at the end of every SUCCESSFUL tick — a wedged tick
     lets it go stale) + state/ALERT (exists while an alarm is active).
     KeepAlive restarts a crashed runner; the pidfile singleton makes a double-start harmless;
     a SIGTERM from `launchctl unload` is a clean fail-stop (nothing merges while it is down). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{label}</string>
    <key>ProgramArguments</key>
    <array>
        <string>{superlooper_bin}</string>
        <string>run</string>
        <string>--repo</string>
        <string>{repo_path}</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{state_home}/logs/runner.log</string>
    <key>StandardErrorPath</key>
    <string>{state_home}/logs/runner.log</string>
</dict>
</plist>
